# systemd/agent-gateway.service
[Unit]
Description=AI agent gateway (hardened, non-root)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Runs as the dedicated non-root agent user — NEVER root.
User=@OC_USER@
Group=@OC_USER@
Environment=HOME=@OC_HOME@
Environment=NODE_ENV=production
Environment=PATH=@AGENT_PATH@
WorkingDirectory=@OC_HOME@
ExecStart=@AGENT_CMD@ gateway
Restart=on-failure
RestartSec=5

# Note: this unit controls the default Hermes profile only. Additional Hermes
# profiles should be run as systemd *system* instances via
# agent-gateway@<profile>.service, not via `hermes --profile <p> gateway start`,
# because the non-login system user has no user D-Bus / linger session.

# Hardening — least privilege for an exec-capable agent on untrusted input.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=@OC_HOME@
ProtectControlGroups=true
ProtectKernelTunables=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
